# crontab
# Apiary
30 * * * * (nice -n 19 python collect-statistics.py)
5 0 * * * (nice -n 19 python collect-general.py)
5 1 * * * (nice -n 19 python collect-extensions.py)
